# packaging/linux/kaitu.service
[Unit]
Description=Kaitu VPN
Documentation=https://kaitu.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/k2 run
Restart=on-failure
RestartSec=3
User=root

# TUN requires CAP_NET_ADMIN; running as root is the simplest path and matches
# the install.sh guidance. If you prefer a dedicated service user, add the
# capability and adjust the User= line accordingly.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN

# Block execve() from picking up setuid/setcap file capabilities on any
# child process (systemctl during self-upgrade, adb helper, etc.). The
# daemon already runs as root so this does not take anything away — it
# only prevents post-exec privilege escalation beyond what the service
# already has.
NoNewPrivileges=yes

# The daemon persists encrypted storage + update state under /etc/kaitu.
# install.sh pre-creates the directory with mode 0700; the daemon also
# calls os.MkdirAll on first write as a fallback.
ReadWritePaths=/etc/kaitu

[Install]
WantedBy=multi-user.target
